# bazarr-autotranslate.timer
[Unit]
Description=Run daily at 2 AM

[Timer]
OnCalendar=*-*-* 9:28:00
Persistent=true

[Install]
WantedBy=timers.target
