# bazarr-autotranslate.service
[Unit]
Description=Auto translate in Bazarr

[Service]
ExecStart=/usr/bin/python3 /path/to/script/bazarr-auto-translate.py
